# roles/edulive/files/livebuilder.service
[Unit]
Description=Run livebuilder script

[Service]
Type=simple
ExecStart=/usr/local/sbin/livebuilder.sh
